# ansible/phases/09-dashboard/roles/dashboard_setup/templates/dashboard-frontend.service.j2
[Unit]
Description=5G Dashboard Frontend (Vite Preview)
After=network-online.target {{ dashboard_backend_service_name }}.service
Wants=network-online.target

[Service]
Type=simple
User={{ dashboard_service_user }}
WorkingDirectory={{ dashboard_frontend_dir }}
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=VITE_FRONTEND_MODE=dev
Environment=VITE_BACKEND_TARGET=http://127.0.0.1:{{ dashboard_backend_port }}
Environment=VITE_AUTH_ENABLED={{ dashboard_auth_enabled | bool | lower }}
Environment=VITE_KEYCLOAK_AUTHORITY={{ dashboard_keycloak_external_url }}{{ dashboard_keycloak_path_prefix }}/realms/{{ dashboard_keycloak_realm }}
Environment=VITE_KEYCLOAK_CLIENT_ID={{ dashboard_keycloak_client_id }}
# HMR over a public tunnel. When dashboard_dev_hmr_enabled is true (the
# default) and a tunnel URL is configured, the browser reaches the dev
# server over wss/443 through that hostname. dashboard_dev_hmr_path lets
# the operator move the WebSocket off the root path so an upstream
# Zero-Trust policy can bypass auth on that single path (root carries
# the SPA auth gate). Set DASHBOARD_DEV_HMR_ENABLED=false to opt out
# entirely (loses auto-reload on save).
{% if dashboard_dev_hmr_enabled | bool and dashboard_dev_external_url %}
Environment=VITE_HMR_HOST={{ dashboard_dev_external_url | regex_replace('^https?://', '') | regex_replace('/.*$', '') }}
Environment=VITE_HMR_PROTOCOL={{ 'wss' if dashboard_dev_external_url.startswith('https://') else 'ws' }}
Environment=VITE_HMR_CLIENT_PORT={{ '443' if dashboard_dev_external_url.startswith('https://') else '80' }}
{% if dashboard_dev_hmr_path %}
Environment=VITE_HMR_PATH={{ dashboard_dev_hmr_path }}
{% endif %}
{% elif not dashboard_dev_hmr_enabled | bool %}
Environment=VITE_HMR_DISABLED=1
{% endif %}
# Unit is only installed when dashboard_dev_enabled=true, so the dev
# server (HMR) is the only mode this service runs.
Environment=CHOKIDAR_USEPOLLING=1
Environment=WATCHPACK_POLLING=true
ExecStart=/usr/bin/npm run dev -- --host 0.0.0.0 --port {{ dashboard_frontend_port }}
Restart=always
RestartSec=3
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
