# ansible/phases/09-dashboard/roles/dashboard_setup/templates/dashboard-frontend.service.j2
[Unit]
Description=5G Dashboard Frontend (Vite Preview)
After=network-online.target {{ dashboard_backend_service_name }}.service
Wants=network-online.target

[Service]
Type=simple
User={{ dashboard_service_user }}
WorkingDirectory={{ dashboard_frontend_dir }}
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=VITE_BACKEND_TARGET=http://127.0.0.1:{{ dashboard_backend_port }}
Environment=VITE_AUTH_ENABLED={{ dashboard_auth_enabled | bool | lower }}
Environment=VITE_KEYCLOAK_PROXY_TARGET={{ dashboard_keycloak_proxy_target }}
Environment=VITE_KEYCLOAK_AUTHORITY={{ dashboard_keycloak_external_url }}{{ dashboard_keycloak_path_prefix }}/realms/{{ dashboard_keycloak_realm }}
Environment=VITE_KEYCLOAK_CLIENT_ID={{ dashboard_keycloak_client_id }}
{% if dashboard_mode == 'dev' %}
Environment=CHOKIDAR_USEPOLLING=1
Environment=WATCHPACK_POLLING=true
ExecStart=/usr/bin/npm run dev -- --host 0.0.0.0 --port {{ dashboard_frontend_port }}
{% else %}
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port {{ dashboard_frontend_port }}
{% endif %}
Restart=always
RestartSec=3
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
