# ansible/phases/01-infrastructure/roles/infra_setup/templates/chrony-force-sync.service.j2
[Unit]
Description=Force chrony time step if drift exceeds threshold

[Service]
Type=oneshot
ExecStart=/usr/bin/chronyc -a makestep
